# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
directory = /etc/nginx
command = /usr/local/sbin/nginx -c /etc/nginx/nginx.conf
autorestart = true

[program:php5-fpm]
directory = /etc/php5
command = /usr/sbin/php5-fpm -R --fpm-config /etc/php5/fpm/php-fpm.conf
autorestart = true


